# doorlock.conf.upstart
# doorlock.conf  - start door lock service

description "door lock webserver and epaper rendering"

start on network


pre-start script
  cd /home/pi/binarykitchen
end script


exec ./webserver.py

respawn
